# dark-alpha-auto-execution/scripts/com.darkalpha.daily-snapshot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.darkalpha.daily-snapshot</string>
  <key>WorkingDirectory</key>
  <string>REPLACE_WITH_ABSOLUTE_PROJECT_PATH</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-l</string>
    <string>-c</string>
    <string>poetry run dark-alpha report daily</string>
  </array>
  <!-- Fire at 00:05 UTC every day.
       macOS uses local time here, so adjust Hour if your system is not UTC. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>0</integer>
    <key>Minute</key><integer>5</integer>
  </dict>
  <key>StandardOutPath</key><string>/tmp/dark-alpha-daily.out</string>
  <key>StandardErrorPath</key><string>/tmp/dark-alpha-daily.err</string>
</dict>
</plist>
